# contrib/systemd/kindling.service
[Unit]
Description=Kindling control plane, builds, and edge proxy
Documentation=https://github.com/Peppermint-Lab/kindling
After=network-online.target
Wants=network-online.target
After=kindling-networking.service
Wants=kindling-networking.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
# buildah/crun need XDG_RUNTIME_DIR. Using /tmp avoids rootless build failures inside user+mount
# namespaces ("mkdir /run/user/<uid>: permission denied" when creating lock dirs).
# Placeholder ubuntu is rewritten by contrib/setup-kindling-prod.sh.
Environment=HOME=/home/ubuntu
Environment=XDG_RUNTIME_DIR=/tmp/kindling-xdg
ExecStartPre=+/bin/sh -c 'install -d -m 0700 -o ubuntu -g ubuntu /tmp/kindling-xdg'
WorkingDirectory=/home/ubuntu/kindling
ExecStart=/usr/local/lib/kindling/serve.sh
Restart=always
RestartSec=5
KillMode=control-group
# Cloud Hypervisor TAP + workloads can need many fds.
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
